# 0-setup/SRB5/etc/profile.d/libintelopencl.sh
export PKG_CONFIG_PATH="/opt/intel/opencl:$PKG_CONFIG_PATH"
